# app/kitchen/site-cookbooks/firehol/templates/default/docker-iptables.conf
description "Start Docker Custom Iptables Rules"
author "ridesharemarket.com"

start on started docker

script

    sleep 3

    /home/ubuntu/docker_iptables.sh

end script
