# app/kitchen/site-cookbooks/firehol/templates/default/etc_init_docker-iptables.conf
description "Start Docker Custom Iptables Rules"

start on started docker

script
    /usr/bin/timeout 30 /home/ubuntu/docker_iptables.sh
end script
