# app/kitchen/site-cookbooks/firehol/templates/default/firehol.conf
description "Start Firehol Iptables Policies"

start on started networking

pre-start exec /sbin/firehol start

post-stop exec /sbin/firehol stop
